# kde/.config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Exec=/home/hgducharme/.config/polybar/launch.sh
Hidden=false
Name=Polybar
