# kde/.config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Exec=/home/hgducharme/.config/polybar/launch.sh --pause=5
Hidden=false
Name=Polybar
